# flexnet/files/profile.d/flexnet.sh
#! /bin/sh
{{ pillar['filebanner']['banner#'] }}
#
export LD_LIBRARY_PATH={{tools_path}}:$LD_LIBRARY_PATH

